# contrib/testnet/systemd/btc-burn-claim-daemon.service
[Unit]
Description=BATHRON BTC Burn Claim Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
ExecStart=/home/ubuntu/BATHRON/contrib/testnet/btc_burn_claim_daemon.sh run
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=btc-burn-claim-daemon

# Environment
Environment=HOME=/home/ubuntu
Environment=BTC_CLI=/home/ubuntu/bitcoin-27.0/bin/bitcoin-cli
Environment=BTC_CONF=/home/ubuntu/.bitcoin-signet/bitcoin.conf
Environment=BATHRON_CLI=/home/ubuntu/bathron-cli

[Install]
WantedBy=multi-user.target
